# etc/supervisord.d/docker-registry.ini
[program:docker-registry]
command=docker-registry /etc/docker-registry/config.yml
redirect_stdery=true
stdout_logfile=/var/log/supervisor/%(program_name)s.log
autorestart=true
autostart=true
startsecs=5
exitcodes=0
